# sonicroom.service
[Unit]
Description=SonicRoom audio conferencing server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/sonicroom
ExecStart=/usr/bin/pnpm start
Restart=on-failure
RestartSec=3

# Load configuration from the repo-root .env (copy it from .env.example).
# Using EnvironmentFile is important: the ANNOUNCED_IP* values are read by
# mediasoup-config.ts at import time, BEFORE the app's own .env loader runs, so
# they must already be in the environment here — otherwise ICE falls back to
# 127.0.0.1 and SFU media (3+ participants) never connects. Put PORT, NODE_ENV,
# ANNOUNCED_IP, ANNOUNCED_IP_LOCAL, INSTANCE_NAME, etc. in that .env.
EnvironmentFile=/home/sonicroom/.env

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/home/sonicroom
ProtectHome=read-only
PrivateTmp=true

[Install]
WantedBy=multi-user.target
